# packages/pipewire-thinkpad-units/pipewire-rtp-to-arch-tower.service
[Unit]
Description=RTP Sender to jonas-tower-arch
After=pipewire.service
BindsTo=pipewire.service

[Service]
Type=simple
Environment=PIPEWIRE_DEBUG=3
ExecStart=%h/.local/bin/pipewire-rtp-to-arch-tower
Restart=on-failure
Slice=session.slice

[Install]
WantedBy=default.target
