# packages/pipewire-tower-units/pipewire-extra-routing.service
[Unit]
Description=Extra Pipewire Routing
After=pipewire.service
BindsTo=pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/pipewire-extra-routing
Restart=on-failure
Slice=session.slice

[Install]
WantedBy=default.target
